# sv/cmdgateway/run
#!/bin/sh

#sv start socklog-unix || exit 1

cd /cmdgateway

exec 2>&1
exec node server.js
