# sv/etcd/run
#!/bin/bash

#sv start socklog-unix || exit 1

ADVERTISE_CLIENT_URLS=$(line=""; for x in $(ip route | grep -Po '(?<=src )(\d{1,3}.){4}'|sort|uniq); do line="http://$x:4001 $line"; done; line=${line%% }; echo ${line// /,})
ETCD_OPTS="-advertise-client-urls=${ADVERTISE_CLIENT_URLS} -listen-client-urls=http://0.0.0.0:4001 -data-dir=/var/lib/etcd-data"

exec 2>&1
exec etcd ${ETCD_OPTS}
