# sv/influxdb/run
#!/bin/bash

#sv start socklog-unix || exit 1

exec 2>&1
exec influxd
