# sv/controller-manager/run
#!/bin/bash

source /etc/envvars

until curl -s http://localhost:8080; do echo "Waiting for API Server..."; sleep 3; done;

export PKI_DIR=/dev/shm/kubernetes

KUBE_CONTROLLER_MANAGER_OPTS="\
--master=http://localhost:8080 \
--root-ca-file=$PKI_DIR/api-server-ca.pem \
--service-account-private-key-file=$PKI_DIR/api-server-key.pem \
--feature-gates=$FEATURE_GATES \
--node-monitor-grace-period=20s \
--pod-eviction-timeout=300s \
--terminated-pod-gc-threshold 100 \
--flex-volume-plugin-dir=/var/lib/kubelet/volumeplugins \
--logtostderr=true \
"

exec 2>&1
exec kube-controller-manager ${KUBE_CONTROLLER_MANAGER_OPTS}
