# sv/apiserver/run
#!/bin/bash

sv start etcd || exit 1

source /etc/envvars

KUBE_APISERVER_OPTS="--insecure-bind-address=0.0.0.0 \
--cors_allowed_origins=.* \
--insecure-port=8080 \
--etcd_servers=http://127.0.0.1:4001 \
--logtostderr=true \
--v=1 \
--service-cluster-ip-range=10.1.0.0/16 \
--runtime-config=extensions/v1beta1/daemonsets=true"

env
exec 2>&1
exec kube-apiserver ${KUBE_APISERVER_OPTS}
