# sv/apiserver/run
#!/bin/bash

source /etc/envvars

export PKI_DIR=/dev/shm/kubernetes
export ROLE=apiserver
export SECURE_PORT=6443

KUBE_APISERVER_OPTS="\
--insecure-port=8080 \
--secure-port=$SECURE_PORT \
--etcd-servers=http://127.0.0.1:4001 \
--service-cluster-ip-range=172.27.0.0/16 \
--feature-gates=$FEATURE_GATES \
--enable-admission-plugins=NamespaceLifecycle,NodeRestriction,LimitRanger,ServiceAccount,DefaultStorageClass,DefaultTolerationSeconds,MutatingAdmissionWebhook,ValidatingAdmissionWebhook,ResourceQuota,PodPreset,Priority \
--authorization-mode=Node,RBAC \
--allow-privileged=true \
--anonymous-auth=false \
--runtime-config="api/all,settings.k8s.io/v1alpha1,scheduling.k8s.io/v1alpha1" \
--client-ca-file=$PKI_DIR/$ROLE-ca.pem \
--tls-cert-file=$PKI_DIR/$ROLE-cert.pem \
--tls-private-key-file=$PKI_DIR/$ROLE-key.pem \
--service-account-key-file=$PKI_DIR/service-account-key.pem \
--requestheader-client-ca-file=$PKI_DIR/$ROLE-ca.pem \
--requestheader-allowed-names= \
--requestheader-extra-headers-prefix=X-Remote-Extra- \
--requestheader-group-headers=X-Remote-Group \
--requestheader-username-headers=X-Remote-User \
--kubelet-client-certificate=$PKI_DIR/$ROLE-cert.pem \
--kubelet-client-key=$PKI_DIR/$ROLE-key.pem \
--proxy-client-cert-file=$PKI_DIR/$ROLE-cert.pem \
--proxy-client-key-file=$PKI_DIR/$ROLE-key.pem \
--enable-aggregator-routing=true \
--storage-backend=etcd3 \
--storage-media-type=application/json \
--event-ttl=1h \
--logtostderr=true \
--v=1 \
"

/vault-init.sh

exec 2>&1
exec kube-apiserver ${KUBE_APISERVER_OPTS}
