# sv/addon-manager/run
#!/bin/bash

source /etc/envvars

exec 2>&1
exec /opt/kube-addons.sh
