# sv/scheduler/run
#!/bin/bash

sv apiserver etcd || exit 1

source /etc/envvars

KUBE_SCHEDULER_OPTS="\
  --logtostderr=true \
  --master=127.0.0.1:8080 \
  --policy-config-file=${SCHEDULER_POLICY_FILE} \
  --v=2"

exec 2>&1
exec kube-scheduler ${KUBE_SCHEDULER_OPTS}
